# build/supervisord-compute.conf
[supervisord]
nodaemon=true
user=root

[program:startup_compute]
priority=1
command=/etc/startup_compute.sh
autorestart=true
