# build/supervisord-ctld.conf
[supervisord]
nodaemon=true

[program:startup_ctld]
priority=1
command=/etc/startup_ctld.sh
autorestart=true
