# build/supervisord-dbd.conf
[supervisord]
nodaemon=true

[program:startup_dbd]
priority=1
command=/etc/startup_dbd.sh
autorestart=true
